# etc/cron.d/open-semantic-search
#
# cron-jobs for Open Semantic Search Apps
#

# RSS Feeds import
* * * * *     opensemanticetl /usr/bin/curl http://apps/search-apps/rss_manager/import >> /dev/null

# Webcrawler
* * * * *     opensemanticetl /usr/bin/curl http://apps/search-apps/crawler/recrawl >> /dev/null

# File crawler
* * * * *     opensemanticetl /usr/bin/curl http://apps/search-apps/files/recrawl >> /dev/null
